# deploy/supervisord.ini
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server
autorestart=true
stdout_logfile=/tmp/redis.log
redirect_stderr=true

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf -g "daemon off;"
autorestart=true
stdout_logfile=/tmp/nginx.log
redirect_stderr=true

[program:daphne]
priority=4
process_name=daphne_%(process_num)s
numprocs=2
directory=/opt/Syntinel/compozed/
command=bash -c "source /opt/venv3/bin/activate && daphne -u /tmp/daphne_%(process_num)s.sock compozed.asgi:channel_layer"
autorestart=true
stdout_logfile=/tmp/daphne_%(process_num)s.log
redirect_stderr=true

[program:worker]
priority=3
process_name=worker_%(process_num)s
command=bash -c "source /opt/venv3/bin/activate && python /opt/Syntinel/compozed/manage.py runworker"
numprocs=4
autorestart=true
stdout_logfile=/tmp/worker_%(process_num)s.log
redirect_stderr=true

[program:loadbalancer]
priority=2
command=/bin/bash -c "cd /opt/Syntinel/distributorCap/src/LoadBalancer && go build . && ./LoadBalancer"
stdout_logfile=/tmp/loadbalancer.log
redirect_stderr=true
autorestart=true

[program:executor1]
priority=1
command=/bin/bash -c "rm -rf dist1 && cp -R dist/ dist1/ && dist1/executor 9090"
stdout_logfile=/tmp/executor1.log
redirect_stderr=true
autorestart=true

[program:executor2]
priority=1
command=/bin/bash -c "rm -rf dist2 && cp -R dist/ dist2/ && dist2/executor 9091"
stdout_logfile=/tmp/executor2.log
redirect_stderr=true
autorestart=true
