# build_image/supervisord.ini
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server
autorestart=true
stdout_logfile=/tmp/redis.log
stderr_logfile=/tmp/redis.error

[program:worker]
process_name=worker_%(process_num)s
command=/usr/bin/python /opt/Syntinel/compozed/manage.py runworker
numprocs=4
autorestart=true
stdout_logfile=/tmp/worker.log
stderr_logfile=/tmp/worker.error

[program:daphne]
process_name=daphne_%(process_num)s
numprocs=2
directory=/opt/Syntinel/compozed/
command=daphne -u /tmp/daphne_%(process_num)s.sock compozed.asgi:channel_layer
autorestart=true
stdout_logfile=/tmp/daphne_%(process_num)s.log
stderr_logfile=/tmp/daphne_%(process_num)s.error

[program:loadbalancer]
command=/bin/bash -c "cd /opt/Syntinel/distributorCap/src/LoadBalancer && go build . && ./LoadBalancer"
stdout_logfile=/tmp/loadbalancer.log
stderr_logfile=/tmp/loadbalancer.error
autorestart=true

[program:executor]
command=/bin/bash -c "cd /opt/Syntinel/syntinel_executor && ./make && dist/executor"
stdout_logfile=/tmp/executor.log
stderr_logfile=/tmp/executor.error
autorestart=true
