# build_image/supervisord.ini
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server
autorestart=true

[program:worker]
process_name=worker_%(process_num)s
command=/usr/bin/python /opt/Syntinel/compozed/manage.py runworker
numprocs=4
autorestart=true

[program:daphne]
directory=/opt/Syntinel/compozed/
command=daphne compozed.asgi:channel_layer
autorestart=true
